# omf/files/omf-resctl.init
#!/bin/sh /etc/rc.common
#
# omf-resctl init script to start omf-resctl daemon
#
START=80


start(){
	omf-resctl &	
}

